# installer.service
[Unit]
Description=Damien-Coopmans-Installer

[Service]
#ExecStartPre=
ExecStart=/root/gns-installer/
SyslogIdentifier=Diskutilization
#ExecStop=

[Install]
WantedBy=multi-user.target
